# docker/php-fpm/build/supervisor/supervisord-programs.conf
[program:app]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/project/bin/console messenger:consume async --memory-limit=512M --time-limit=3600
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
user=project
startsecs=1
startretries=10
stdout_logfile=/var/www/project/var/log/worker.log
